# crontab
0 20 2 * * * * ./main.py # every day at 02:00
